# system_files/etc/systemd/system/himmelblau-tasks.service.d/10-selinux-install.conf
[Unit]
After=himmelblau-selinux-install.service
Requires=himmelblau-selinux-install.service